# system_monitoring/rootfs/etc/services.d/system_monitor/run
#!/usr/bin/with-contenv bashio

# 1. Variablen definieren (Zuerst!)
BLUE='\e[1;34m'
GREEN='\e[1;32m'
NC='\e[0m'
BOLD='\e[1m'

# 2. Header ausgeben
echo -e "${BLUE}╔═══════════════════════════════════════════════════════════╗${NC}"
echo -e "${BLUE}║  🚀 HOME ASSISTANT SYSTEM MONITORING                        ║${NC}"
echo -e "${BLUE}╚═══════════════════════════════════════════════════════════╝${NC}"
echo -e " ${BLUE}▸${NC} ${BOLD}STATUS:${NC}    Initializing Core Components..."

# 3. Konfiguration für Node.js exportieren
export MQTT_BROKER=$(bashio::config 'mqtt_broker')
export MQTT_PORT=$(bashio::config 'mqtt_port')
export MQTT_USER=$(bashio::config 'mqtt_user')
export MQTT_PASSWORD=$(bashio::config 'mqtt_password')
export MQTT_TOPIC_PREFIX=$(bashio::config 'mqtt_topic_prefix')
export UPDATE_INTERVAL=$(bashio::config 'update_interval')
export ENABLE_MQTT=$(bashio::config 'enable_mqtt')
export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"

# 4. Node.js starten
cd /app
exec node server.js
